# ElecSim/example.cron
59 23 * * * /home/pieter/Desktop/DB/ElecSim/run.bash
